# scripts/nous-proxy.service
[Unit]
Description=NousResearch OAuth Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nous-proxy
Group=nous-proxy
WorkingDirectory=/opt/nous-proxy
EnvironmentFile=/opt/nous-proxy/.env
ExecStart=/opt/nous-proxy/.venv/bin/nous-proxy
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/nous-proxy/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
